# write.service

[Unit]
Description=Write time into driver

[Service]
ExecStart =/bin/bash '/home/ubuntu/test/write.sh'

