# write.service

[Unit]
Description=Write time into driver

[Service]
ExecStart =/bin/bash '/tmp/write.sh'

