# write.timer
[Unit]
Description=test
[Timer]
OnUnitActiveSec=30s
[Install]
WantedBy=timers.target
