# write.timer
[Unit]
Description=test
[Timer]
OnUnitActiveSec=120s
[Install]
WantedBy=timers.target
